# turtle-kiosk-simple.service
[Unit]
Description=Turtle Dashboard Browser (Firefox)
After=graphical-session.target turtle-dashboard.service
Wants=graphical-session.target

[Service]
Type=simple
User=shrimp
Group=shrimp
Environment=DISPLAY=:0
ExecStart=/snap/bin/firefox http://localhost:8080
Restart=always
RestartSec=10

[Install]
WantedBy=graphical-session.target 