# hardware/temperhum-manager.service
[Unit]
Description=TEMPerHUM USB Sensor Manager
Documentation=https://github.com/jopey-woof/turtx
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=shrimp
Group=shrimp
WorkingDirectory=/home/shrimp/turtle-monitor
ExecStart=/usr/bin/python3 /home/shrimp/turtle-monitor/hardware/temperhum_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /tmp /home/shrimp/turtle-monitor

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target 