# hardware/temperhum-capture.service
[Unit]
Description=TemperhUM Sensor Data Capture
After=network.target mqtt.service
Wants=mqtt.service

[Service]
Type=simple
User=shrimp
Group=shrimp
WorkingDirectory=/home/shrimp/turtle-monitor/hardware
ExecStart=/usr/bin/python3 /home/shrimp/turtle-monitor/hardware/simple_data_capture.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /var/log /home/shrimp/turtle-monitor

[Install]
WantedBy=multi-user.target 