# hardware/temperhum-manager.service
[Unit]
Description=TEMPerHUM USB Sensor Manager
Documentation=https://github.com/jopey-woof/turtx
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=shrimp
Group=shrimp
ExecStart=/usr/bin/python3 /home/shrimp/turtle-monitor/hardware/temperhum_fixed.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0

# Security settings (simplified for compatibility)
ReadWritePaths=/var/log /tmp

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target 