# turtle-monitor/kiosk/turtle-monitor-kiosk.service
[Unit]
Description=Turtle Monitor Kiosk Display
After=network.target docker.service
Wants=docker.service
Requires=network.target

[Service]
Type=simple
User=shrimp
Group=shrimp
WorkingDirectory=/home/shrimp/turtx
Environment=DISPLAY=:0
Environment=HOME=/home/shrimp
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_SESSION_TYPE=x11
Environment=XAUTHORITY=/home/shrimp/.Xauthority

# Wait for X server to be ready
ExecStartPre=/bin/bash -c "until pgrep Xorg; do sleep 2; done"
ExecStartPre=/bin/bash -c "sleep 5"

# Wait for Turtle Monitor API to be accessible
ExecStartPre=/bin/bash -c "while ! curl -s http://localhost:8000/api/health > /dev/null; do sleep 5; done"

ExecStart=/home/shrimp/turtx/kiosk/start-turtle-monitor-kiosk.sh
Restart=always
RestartSec=15
StartLimitInterval=0
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Don't kill the service if it exits normally
RemainAfterExit=no

[Install]
WantedBy=multi-user.target 